# gui/autostart/protonvpn.desktop
[Desktop Entry]
Type=Application
Name=ProtonVPN
Exec=protonvpn-app --start-minimized
Hidden=false